# deploy/crusade-deck.service
[Unit]
Description=Crusade Deck game server (Colyseus)
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/crusade-deck/server
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=2

Environment=NODE_ENV=production
Environment=PORT=2567
# Сервер слушает 0.0.0.0, наружу его не пускает UFW: открыты только 22/80/443,
# 2567 доступен лишь Caddy по петле. См. deploy/README.md, шаг 1.

# Аккаунты пишутся в server/data — единственный каталог, куда нужна запись.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/deploy/crusade-deck/server/data

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
